# conf/supervisord.conf
[supervisord]
nodaemon = true

[program:httpd]
# TODO: Why does this fail as apache user?
# user = apache
command = /bin/bash -c "exec /sbin/httpd -D FOREGROUND"
